# src/scripts/folding-s-symbiote.sentinel.service
[Unit]
Description=folding.S.symbiote sentinel service
After=network.target

[Service]
Type=simple
ExecStart=python /usr/bin/folding-s-symbiote.py

[Install]
WantedBy=multi-user.target 